# infrastructure/launchd/com.lrghomes.crms-reply-detection.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily reply detection for the Relationships tab. For every sent
  relationship_touch, checks the contact's iMessage history (chat.db, via the
  CRMS sidecar's /sync-imessage) for an inbound reply within 7 days and writes
  relationship_touches.replied_at. Brief: briefs/RELATIONSHIPS_REPLY_TRACKING.md

  Install (Mac mini):
    cp infrastructure/launchd/com.lrghomes.crms-reply-detection.plist \
       ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.lrghomes.crms-reply-detection.plist

  Trigger once manually:
    launchctl start com.lrghomes.crms-reply-detection

  Logs:
    /tmp/lrg-crms-reply-detection.log
    /tmp/lrg-crms-reply-detection-err.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lrghomes.crms-reply-detection</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>scripts/detect-relationship-replies.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ryanlarocca/.openclaw/workspace/PROJECTS/mission-control</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/lrg-crms-reply-detection.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/lrg-crms-reply-detection-err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
